# runtimes/php/8.4/crontab
# /etc/crontabs/ for laravel
# min   hour    day     month   weekday command
*       *       *       *       *       cd /var/www/html && php artisan schedule:run >> /dev/null 2>&1